# ansible/roles/system-bootstrap/templates/timezone-dhaka.service.j2
[Unit]
Description=Set the time zone
[Service]
ExecStart=/usr/bin/timedatectl set-timezone {{ COREOS_TIMEZONE }}
RemainAfterExit=yes
Type=oneshot
