# ansible/roles/import-users/files/import-users.service
[Unit]
Description=Import user from AWS IAM

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c '/opt/bin/import_users.sh'
