# ansible/roles/essential-container/templates/sweeper.timer.j2
[Unit]
Description=Run sweeper dailly

[Timer]
OnBootSec=600
OnCalendar=daily

[Install]
WantedBy=multi-user.target
